# supervisor.conf
[program:daphne]
command=bash /home/www/projects/avitolog/start_django_service.sh
user=www
process_name=daphne
numproc=1
autostart=1
autorestart=1
redirect_strerr=true